# etc/barnacles-backup.timer
[Unit]
Description=barnacles backup dot timer

[Timer]
OnBootSec=12h
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
